# linxian-fresh-miniapp/deploy/systemd.service.example
[Unit]
Description=Linxian Fresh Miniapp Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/linxian-fresh-miniapp/server
ExecStart=/usr/bin/node app.js
Restart=always
Environment=NODE_ENV=production
Environment=APP_PORT=3000
User=www-data

[Install]
WantedBy=multi-user.target
